# updater_onboard/src/deb-pkg-files/postinst
#!/bin/sh

PKG_NAME=oresat-linux-updater

# load new daemon service file and start
systemctl daemon-reload
systemctl enable $PKG_NAME

